# continuousop/deploy/ansible/templates/supervisor.ini
[program:{{ app_name }}]
command=/usr/local/bin/gunicorn --config {{ webapps_dir }}/{{ app_name }}/gunicorn.conf.py wsgi
directory={{ webapps_dir }}/{{ app_name }}/django
autostart=true
autorestart=true
stopsignal=KILL
killasgroup=true
