# irc_server_deployment/hubot/templates/supervisor.ini
[program:hubot]
command={{ hubot_directory }}/node_modules/.bin/hubot -a irc
directory={{ hubot_directory }}
autostart=true
autorestart=true
stopsignal=KILL
killasgroup=true
environment=HUBOT_IRC_SERVER="{{ irc_server_name }}",HUBOT_IRC_NICK="{{ hubot_name }}",HUBOT_IRC_ROOMS="{{ irc_channel_name }}",PORT={{ irc_server_port }}

