# continuousop/deploy/supervisor/continuousop.ini
[program:{{ app_name }}]
command=/usr/bin/python manage.py runserver
directory={{ webapps_dir }}/{{ app_name }}/django
autostart=true
autorestart=true
stopsignal=KILL
killasgroup=true
