# continuousop/deploy/ansible/hubot/templates/supervisor.ini
[program:hubot]
command={{ hubot_directory }}/node_modules/.bin/hubot -a irc
directory={{ hubot_directory }}
autostart=true
autorestart=true
stopsignal=KILL
killasgroup=true
environment=HUBOT_IRC_SERVER="babylon1.willowgarage.com",HUBOT_IRC_NICK="hubot",HUBOT_IRC_ROOMS="#im",PORT=6667

